# conf/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=0 ; (num of main logfile rotation backups;default 10)
loglevel=info ; (log level;default info; others: debug,warn,trace)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:monkey]
command=/monkey/monkey serve
autostart=true
autorestart=true
startretries=10
user=root
directory=/monkey
redirect_stderr=true
stdout_logfile=/var/log/mathilde/monkey.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=4

[eventlistener:monkey_exit]
command=python /monkey/kill_supervisor.py
process_name=monkey
events=PROCESS_STATE_FATAL